# attic/benjamin/var/benjamintoll.com/benjamin.cron
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/CRONGbFmrR installed on Thu Jun 15 16:35:07 2006)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
SHELL=/bin/sh
PATH=/sbin:/usr/sbin:/bin:/usr/bin:/usr/X11R6/bin
MAILTO="ben@benjamintoll.com"
0 0 * * * lynx http://www.benjamintoll.com/italy/lib/send_wotd.php
